# server/deploy/robot-islands-auth.service
[Unit]
# NOTE: named *-auth (not *-server) on purpose — robot-islands-server.service is
# already taken by the separate /root/islands rewrite (ws on :8090). Do not rename
# this back without checking that collision.
Description=Robot Islands auth server (migration slice 1; HTTP on 127.0.0.1:5180)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
WorkingDirectory=/root/robot-islands/server
EnvironmentFile=/root/robot-islands/server/.env
# systemd has no $USER/$PGUSER; the pure-JS pg driver needs a user for the
# postgresql:///db socket form. Supply it explicitly (peer auth -> role root).
Environment=PGUSER=root
# /usr/bin/node on this box is v12; Fastify 5 needs >=20. Use the nvm node
# (same as robot-islands-dev.service) and put it on PATH for child processes.
Environment=PATH=/root/.nvm/versions/node/v22.22.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/root/.nvm/versions/node/v22.22.0/bin/node --import tsx /root/robot-islands/server/src/index.ts
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
